# packaging/systemd/medium-home-node.service
[Unit]
Description=Medium Home Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=OVERLAY_CONTROL_URL={{CONTROL_URL}}
Environment=OVERLAY_SHARED_SECRET={{SHARED_SECRET}}
ExecStart={{HOME_NODE_BIN}} --config {{NODE_CONFIG_PATH}}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
