# packaging/systemd/medium-control-plane.service
[Unit]
Description=Medium Control Plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{STATE_DIR}}
Environment=MEDIUM_CONTROL_DATABASE_URL={{DATABASE_URL}}
ExecStart={{MEDIUM_BIN}} serve control-plane --config {{CONTROL_CONFIG_PATH}}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
