# packaging/systemd/medium-node-agent.service
[Unit]
Description=Medium Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=OVERLAY_CONTROL_URL={{CONTROL_URL}}
Environment=OVERLAY_SHARED_SECRET={{SHARED_SECRET}}
Environment=MEDIUM_CONTROL_PIN={{CONTROL_PIN}}
Environment=MEDIUM_RELAY_ADDR={{RELAY_ADDR}}
{{WSS_RELAY_ENV}}
ExecStart={{NODE_AGENT_BIN}} --config {{NODE_CONFIG_PATH}}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
